# files/etc/crontabs/root
10 4 * * * sleep 5 && touch /etc/banner && reboot
0 1 * * * echo "" > /tmp/ssrplus.log
0 5 * * * /usr/share/shadowsocksr/update.sh
